# install/gravis-gunicorn.socket
[Unit]
Description=GRAVIS gunicorn socket
PartOf=gravis.service

[Socket]
ListenStream=/run/gravis-gunicorn.sock

[Install]
WantedBy=sockets.target
